# app/supervisor/conf.available.d/workers.conf
[program:cms_low_3]

environment=CONCURRENCY=3,LOGLEVEL=info,DJANGO_SETTINGS_MODULE=aws,PYTHONPATH=/edx/app/edxapp/edx-platform,SERVICE_VARIANT=cms
user=www-data
directory=/edx/app/edxapp/edx-platform
stdout_logfile=/edx/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/edx/var/log/supervisor/%(program_name)s-stderr.log

command=/edx/app/edxapp/venvs/edxapp/bin/python /edx/app/edxapp/edx-platform/manage.py cms --settings=aws celery worker --loglevel=info --queues=edx.cms.core.low --hostname=edx.cms.core.low.%%h --concurrency=3 
killasgroup=true
stopasgroup=true
stopwaitsecs=432000

[program:cms_default_4]

environment=CONCURRENCY=4,LOGLEVEL=info,DJANGO_SETTINGS_MODULE=aws,PYTHONPATH=/edx/app/edxapp/edx-platform,SERVICE_VARIANT=cms
user=www-data
directory=/edx/app/edxapp/edx-platform
stdout_logfile=/edx/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/edx/var/log/supervisor/%(program_name)s-stderr.log

command=/edx/app/edxapp/venvs/edxapp/bin/python /edx/app/edxapp/edx-platform/manage.py cms --settings=aws celery worker --loglevel=info --queues=edx.cms.core.default --hostname=edx.cms.core.default.%%h --concurrency=4 
killasgroup=true
stopasgroup=true
stopwaitsecs=432000

[program:cms_high_1]

environment=CONCURRENCY=1,LOGLEVEL=info,DJANGO_SETTINGS_MODULE=aws,PYTHONPATH=/edx/app/edxapp/edx-platform,SERVICE_VARIANT=cms
user=www-data
directory=/edx/app/edxapp/edx-platform
stdout_logfile=/edx/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/edx/var/log/supervisor/%(program_name)s-stderr.log

command=/edx/app/edxapp/venvs/edxapp/bin/python /edx/app/edxapp/edx-platform/manage.py cms --settings=aws celery worker --loglevel=info --queues=edx.cms.core.high --hostname=edx.cms.core.high.%%h --concurrency=1 
killasgroup=true
stopasgroup=true
stopwaitsecs=432000

[program:lms_low_1]

environment=CONCURRENCY=1,LOGLEVEL=info,DJANGO_SETTINGS_MODULE=aws,PYTHONPATH=/edx/app/edxapp/edx-platform,SERVICE_VARIANT=lms
user=www-data
directory=/edx/app/edxapp/edx-platform
stdout_logfile=/edx/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/edx/var/log/supervisor/%(program_name)s-stderr.log

command=/edx/app/edxapp/venvs/edxapp/bin/python /edx/app/edxapp/edx-platform/manage.py lms --settings=aws celery worker --loglevel=info --queues=edx.lms.core.low --hostname=edx.lms.core.low.%%h --concurrency=1 
killasgroup=true
stopasgroup=true
stopwaitsecs=432000

[program:lms_default_3]

environment=CONCURRENCY=3,LOGLEVEL=info,DJANGO_SETTINGS_MODULE=aws,PYTHONPATH=/edx/app/edxapp/edx-platform,SERVICE_VARIANT=lms
user=www-data
directory=/edx/app/edxapp/edx-platform
stdout_logfile=/edx/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/edx/var/log/supervisor/%(program_name)s-stderr.log

command=/edx/app/edxapp/venvs/edxapp/bin/python /edx/app/edxapp/edx-platform/manage.py lms --settings=aws celery worker --loglevel=info --queues=edx.lms.core.default --hostname=edx.lms.core.default.%%h --concurrency=3 
killasgroup=true
stopasgroup=true
stopwaitsecs=432000

[program:lms_high_4]

environment=CONCURRENCY=4,LOGLEVEL=info,DJANGO_SETTINGS_MODULE=aws,PYTHONPATH=/edx/app/edxapp/edx-platform,SERVICE_VARIANT=lms
user=www-data
directory=/edx/app/edxapp/edx-platform
stdout_logfile=/edx/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/edx/var/log/supervisor/%(program_name)s-stderr.log

command=/edx/app/edxapp/venvs/edxapp/bin/python /edx/app/edxapp/edx-platform/manage.py lms --settings=aws celery worker --loglevel=info --queues=edx.lms.core.high --hostname=edx.lms.core.high.%%h --concurrency=4 
killasgroup=true
stopasgroup=true
stopwaitsecs=432000

[program:lms_high_mem_2]

environment=CONCURRENCY=2,LOGLEVEL=info,DJANGO_SETTINGS_MODULE=aws,PYTHONPATH=/edx/app/edxapp/edx-platform,SERVICE_VARIANT=lms
user=www-data
directory=/edx/app/edxapp/edx-platform
stdout_logfile=/edx/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/edx/var/log/supervisor/%(program_name)s-stderr.log

command=/edx/app/edxapp/venvs/edxapp/bin/python /edx/app/edxapp/edx-platform/manage.py lms --settings=aws celery worker --loglevel=info --queues=edx.lms.core.high_mem --hostname=edx.lms.core.high_mem.%%h --concurrency=2 --maxtasksperchild 1
killasgroup=true
stopasgroup=true
stopwaitsecs=432000


[group:edxapp_worker]
programs=cms_low_3,cms_default_4,cms_high_1,lms_low_1,lms_default_3,lms_high_4,lms_high_mem_2
